# deploy/macos/com.careerworkflow.ui.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
<key>Label</key><string>com.careerworkflow.ui</string>
<key>ProgramArguments</key><array><string>REPLACE_PYTHON</string><string>REPLACE_REPO/run_nicegui.py</string></array>
<key>WorkingDirectory</key><string>REPLACE_REPO</string>
<key>RunAtLoad</key><true/><key>KeepAlive</key><true/>
<key>StandardOutPath</key><string>REPLACE_REPO/logs/ui.out.log</string>
<key>StandardErrorPath</key><string>REPLACE_REPO/logs/ui.err.log</string>
</dict></plist>
